# cookbooks/development/files/default/etc/profile.d/grails.sh
GRAILS_HOME=/usr/local/grails/grails-2.0.3
PATH=/usr/bin/grailsx:$PATH